# examples/dashboard-server.service
[Unit]
Description=Status Dashboard Server

[Service]
ExecStart=/usr/local/bin/dashboard-server -c /etc/status-dashboard/server.json
User=dashboard
